# music-widget.service
[Unit]
Description=Apple Music Now Playing widget
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# %h expands to the user's home directory
WorkingDirectory=%h/apple-music-now-playing
ExecStart=/usr/bin/python3 %h/apple-music-now-playing/start.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
